# udev/99-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", RUN+="/bin/bash /home/raabls/.screenlayout/screenlayout.sh udev"
