# udev/99-monitor.rules
ACTION=="change", KERNEL=="card0", SUBSYSTEM=="drm", RUN+="/bin/bash /home/raabls/.screenlayout/screenlayout.sh udev"
